# lain/.config/systemd/user/compton.service
[Unit]
Description=Compton x11 compositor

[Service]
Type=forking
ExecStart=compton --config %E/compton/config -b

[Install]
WantedBy=i3@.service
