# lain/.config/systemd/user/ibus.service
[Unit]
Description=IBus input daemon

[Service]
Type=simple
ExecStart=ibus-daemon --xim

[Install]
WantedBy=i3@.service
