# lain/.config/systemd/user/polybar-top.service
[Unit]
Description=Polybar Top

[Service]
Type=simple
ExecStart=polybar top

[Install]
WantedBy=i3@.service
