# lain/.config/systemd/user/polybar-bottom.service
[Unit]
Description=Polybar Bottom

[Service]
Type=simple
ExecStart=polybar bottom

[Install]
WantedBy=i3@.target
