# lain/.config/systemd/user/wallpaper.timer
[Unit]
Description=Periodically run wallpaper service

[Timer]
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
